# packages/lan/tamanu-facility.service
[Unit]
Description=Tamanu Facility Server
After=network.target

[Service]
ExecStart=/usr/bin/node dist
WorkingDirectory=/var/lib/tamanu/packages/lan
User=tamanu
Group=tamanu

CPUAccounting=yes
MemoryAccounting=yes
TasksAccounting=yes
IOAccounting=yes
IPAccounting=yes

[Install]
WantedBy=multi-user.target